# server/dist/systemd.service
[Unit]
Description=capitaladay app
After=network.target

[Service]
EnvironmentFile=/usr/share/capitaladay/env
PIDFile=/var/run/capitaladay/RUNNING_PID
WorkingDirectory=/
ExecStart=/usr/bin/capitaladay -Dhttp.port=9001 -Dpidfile.path=/var/run/capitaladay/RUNNING_PID -Dplay.evolutions.db.default.autoApply=true
Restart=on-failure
User=capitaladay
Group=capitaladay

[Install]
WantedBy=multi-user.target
